# plist/com.mindvault.arctic-ko-mlx.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.mindvault.arctic-ko-mlx</string>

    <key>ProgramArguments</key>
    <array>
        <string>__USER_HOME__/.claude/scripts/mindvault/arctic_ko_server_runner.sh</string>
    </array>

    <key>EnvironmentVariables</key>
    <dict>
        <key>HOME</key>
        <string>__USER_HOME__</string>
        <key>HF_HOME</key>
        <string>__USER_HOME__/.cache/huggingface</string>
        <key>PYTHONUNBUFFERED</key>
        <string>1</string>
    </dict>

    <key>WorkingDirectory</key>
    <string>__USER_HOME__</string>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key>
        <false/>
        <key>Crashed</key>
        <true/>
    </dict>

    <key>ThrottleInterval</key>
    <integer>30</integer>

    <key>StandardOutPath</key>
    <string>__USER_HOME__/Library/Logs/arctic-ko-mlx.log</string>

    <key>StandardErrorPath</key>
    <string>__USER_HOME__/Library/Logs/arctic-ko-mlx.err</string>

    <key>ProcessType</key>
    <string>Interactive</string>
</dict>
</plist>
